# roles/injectd/files/injectd.service
[Unit]
Description=Inject

[Service]
ExecStart=/usr/bin/injectd

[Install]
WantedBy=multi-user.target